# wallepy.timer
[Unit]
Description=Run wallepy every hour and on boot

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
